# images/shaddock/qemu/supervisord.conf
[supervisord]
nodaemon=true

[program:libvirtd]
command=/usr/sbin/libvirtd
autorestart=true
priority=20

[program:nova-compute]
command=/usr/bin/nova-compute
autorestart=true
priority=30
