# template/keystone/supervisord.conf
[supervisord]
nodaemon=true

[program:keystone-all]
command=/usr/bin/keystone-all
autorestart=true
priority=30
