# images/ubuntu-mitaka/shaddock/keystone-git/supervisord.conf
[supervisord]
nodaemon=true

[program:keystone-all]
command=/opt/keystone/bin/keystone-all --config-file /opt/keystone/etc/keystone.conf
autorestart=true
priority=30
