# template/heat/supervisord.conf
[supervisord]
nodaemon=true

[program:neutron-all]
command=/usr/bin/neutron-all
autorestart=true
priority=30
