# template/cinder/supervisord.conf
[supervisord]
nodaemon=true

[program:cinder-all]
command=/usr/bin/cinder-all
autorestart=true
priority=30
