# images/shaddock/block-storage/supervisord.conf
[supervisord]
nodaemon=true

[program:tgtd]
command=/usr/sbin/tgtd -f
autorestart=true
priority=20

[program:cinder-volume]
command=/usr/bin/cinder-volume
autorestart=true
priority=30
