# images/shaddock/block-storage/supervisord.conf
[supervisord]
nodaemon=true

[program:cinder-registry]
command=/usr/sbin/tgtd -f
autorestart=true
priority=20

[program:cinder-api]
command=/usr/bin/cinder-volume
autorestart=true
priority=30
