# images/region1/shaddock/nova/supervisord.conf
[supervisord]
nodaemon=true

[program:nova-api]
command=/opt/openstack/services/nova/bin/nova-api
autorestart=true
priority=30

[program:nova-cert]
command=/opt/openstack/services/nova/bin/nova-cert
autorestart=true
priority=30

[program:nova-consoleauth]
command=/opt/openstack/services/nova/bin/nova-consoleauth
autorestart=true
priority=30

[program:nova-scheduler]
command=/opt/openstack/services/nova/bin/nova-scheduler
autorestart=true
priority=30

[program:nova-conductor]
command=/opt/openstack/services/nova/bin/nova-conductor
autorestart=true
priority=30

[program:nova-novncproxy]
command=/opt/openstack/services/nova/bin/nova-novncproxy
autorestart=true
priority=30
