# images/ubuntu-openstack-newton/shaddock/kvm/supervisord.conf
[supervisord]
nodaemon=true

[program:kvm]
command=/usr/local/bin/kvm.sh
autorestart=true
priority=20

[program:nova-compute]
command=/usr/bin/nova-compute
autorestart=true
priority=30