# images/shaddock/neutron/supervisord.conf
[supervisord]
nodaemon=true

[program:neutron-server]
command=/usr/bin/neutron-server
autorestart=true
priority=30

[program:neutron-linuxbridge-agent]
command=/usr/bin/neutron-linuxbridge-agent
autorestart=true
priority=30

[program:neutron-dhcp-agent]
command=/usr/bin/neutron-dhcp-agent
autorestart=true
priority=30

[program:neutron-metadata-agent]
command=/usr/bin/neutron-metadata-agent
autorestart=true
priority=30

[program:neutron-l3-agent]
command=/usr/bin/neutron-l3-agent
autorestart=true
priority=30
