# images/shaddock/cinder-volume/supervisord.conf
[supervisord]
nodaemon=true

[program:tgt]
command=/usr/sbin/tgtd -f
autorestart=true
priority=30

[program:cinder-volume]
command=/usr/bin/cinder-volume
autorestart=true
priority=30
