# images/shaddock/manila-share/supervisord.conf
[supervisord]
nodaemon=true

[program:manila-scheduler]
command=/usr/bin/manila-share --log-file=/var/log/manila/manila-share.log
autorestart=true
priority=30
