# images/shaddock/cinder/supervisord.conf
[supervisord]
nodaemon=true

[program:cinder-registry]
command=/usr/bin/cinder-registry
autorestart=true
priority=30

[program:cinder-api]
command=/usr/bin/cinder-api
autorestart=true
priority=30
